# scripts/totp-server.service
[Unit]
Description=TOTP Authentication Server
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=totp
Group=totp
WorkingDirectory=/opt/totp/server
EnvironmentFile=/opt/totp/.env
ExecStart=/usr/bin/node src/app.js
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/totp/server.log
StandardError=append:/var/log/totp/server-error.log
SyslogIdentifier=totp-server

# 资源限制
LimitNOFILE=65536
MemoryMax=512M

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/totp/uploads /opt/totp/logs /var/log/totp

[Install]
WantedBy=multi-user.target
