# scripts/totp-webadmin.service
[Unit]
Description=TOTP Web Admin
After=network.target totp-server.service
Wants=totp-server.service

[Service]
Type=simple
User=totp
Group=totp
WorkingDirectory=/opt/totp/web-admin
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 5174
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/totp/web-admin.log
StandardError=append:/var/log/totp/web-admin-error.log

[Install]
WantedBy=multi-user.target
